# linux/.config/autostart/ssh-agent-session-sync.desktop
[Desktop Entry]
Version=1.0
Type=Application
Name=SSH Agent Session Sync
Comment=Start or reconnect to the shared SSH agent for desktop apps
Exec=/home/metsatron/.local/bin/ssh-agent-session-sync
Terminal=false
StartupNotify=false
X-GNOME-Autostart-enabled=true
