# x230/.config/systemd/user/kdeconnect-refresh.service
[Unit]
Description=Repair KDE Connect (refresh + wait + act)

[Service]
Type=oneshot
ExecStart=/home/metsatron/.local/bin/kdeconnect-repair
