# think/.config/systemd/user/kdeconnect_refresh.service
[Unit]
Description=Refresh KDE Connect if no devices are available

[Service]
Type=oneshot
# Only run when nothing is "available" (paired + reachable)
ExecCondition=/bin/bash -lc '! kdeconnect-cli --list-available | grep -q .'
# Refresh, then poll up to ~10s for any available device
ExecStart=/bin/bash -lc 'kdeconnect-cli --refresh; for i in {1..20}; do kdeconnect-cli --list-available | grep -q . && exit 0; sleep 0.5; done; exit 0'